# after-37-backup-checker-1/ansible/certbot/files/certbot.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/root/certbot/get-certs.sh
